# unraid/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

; ── PostgreSQL ─────────────────────────────────────────────────────────────────
[program:postgres]
command=/usr/local/bin/init-db.sh
priority=10
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── Redis ──────────────────────────────────────────────────────────────────────
[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --save "" --appendonly no --loglevel notice
priority=20
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── Backend API ────────────────────────────────────────────────────────────────
[program:backend]
command=/usr/local/bin/start-backend.sh
priority=30
autorestart=true
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── Webhook worker ─────────────────────────────────────────────────────────────
[program:webhook-worker]
command=/usr/local/bin/start-worker.sh
priority=40
autorestart=true
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── Next.js frontend ───────────────────────────────────────────────────────────
[program:frontend]
command=/usr/local/bin/start-frontend.sh
priority=50
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
